# templates/tinc.systemd.service.j2
[Unit]
Description=tinc vpn
After=network.target

[Service]
Type=forking
ExecStart={{ tinc_systemd_execstart }}
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
